# docker/dev/supervisor/app.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true

[program:rails]
command=bundle exec pumactl start
directory=/var/www/app
